# systemd-unit/credfeto-podfox.timer
[Unit]
Description=Update Podfox

[Timer]
OnBootSec=5min
OnCalendar=*:0/5
AccuracySec=30

[Install]
WantedBy=multi-user.target
